# ops/monitoring/uptime-monitor.timer
[Unit]
Description=Start uptime monitor at boot

[Timer]
Unit=uptime-monitor.service
OnBootSec=30
Persistent=true

[Install]
WantedBy=timers.target
