# packaging/systemd/cloudflared-healthcheck.timer
[Unit]
Description=Run cloudflared healthcheck every 30s

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=1s

[Install]
WantedBy=timers.target
